# s6/smokeping_prober/run
#!/bin/sh
set -eu

exec /usr/bin/smokeping_prober
